# guest-portal.service
[Unit]
Description=Guest Portal Node.js Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/guest-portal
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
