# guest-portal.service
[Unit]
Description=Guest Portal Node.js Server
After=network.target

[Service]
Type=simple
User=guestportal
Group=guestportal
WorkingDirectory=/opt/guest-portal
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/etc/guest-portal /opt/guest-portal/uploads /mnt
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
